# database/backup-service/crontab
0 1 * * * /backup.sh >> /var/log/backup.log